# src/sra/sdk/autobuild/mac32.obelix.crontab
* * * * * /bin/bash /net/pan1/sra-test/autobuild/32/autobuild/build.sh M32 /net/pan1/sra-test/autobuild/32/autobuild /net/pan1/sra-test/autobuild/32/ running static
